# example.timer
[Unit]
Description=Example

[Timer]
OnCalendar=*-*-* 15:00:00
Persistent=true
Unit=example.service

[Install]
WantedBy=multi-user.target
